# sql/crontab
#rollup jobs
* * * * * /opt/pgsql835/bin/psql -d reconnoiter -U stratcon -c "select stratcon.rollup_matrix_numeric_5m();" >/tmp/rollup_5.log 2>&1
00,10,20,30,40,50 * * * * /opt/pgsql835/bin/psql -d reconnoiter -U stratcon -c "select stratcon.rollup_matrix_numeric_20m();" >/tmp/rollup_20.log 2>&1
01,31 * * * * /opt/pgsql835/bin/psql -d reconnoiter -U stratcon -c "select stratcon.rollup_matrix_numeric_60m();" >/tmp/rollup_60.log 2>&1
02 00,03,06,09,12,15,18 * * * /opt/pgsql835/bin/psql -d reconnoiter -U stratcon -c "select stratcon.rollup_matrix_numeric_6hours();" >/tmp/rollup_6hours.log 2>&1
03 00,06,12 * * * /opt/pgsql835/bin/psql -d reconnoiter -U stratcon -c "select stratcon.rollup_matrix_numeric_12hours();" >/tmp/rollup_12hours.log 2>&1
#cleanup jobs
01 00 * * * /opt/pgsql835/bin/psql -d reconnoiter -c "delete from prism.saved_graphs where saved = false and last_update<current_timestamp - '1 day' ::interval;" 1>/dev/null